# root-rutorrent/etc/services.d/php-fpm7/run
#!/usr/bin/execlineb -P
with-contenv
importas -s opts PHP7_OPTS
foreground { echo ${opts} }
php-fpm7 -F ${opts}