# root/etc/services.d/httpd/run
#!/usr/bin/execlineb -P
with-contenv

# Apache gets grumpy about PID files pre-existing
foreground { rm -f /usr/local/apache2/logs/httpd.pid }
/usr/local/apache2/bin/httpd -DFOREGROUND