# root/etc/services.d/rtorrent/run
#!/usr/bin/execlineb -P

with-contenv
s6-setuidgid abc
importas -s opts RTORRENT_OPTS

# importas -D"-n -o system.daemon.set=true" -s opts RTORRENT_OPTS
# /usr/bin/rtorrent ${opts}

foreground { echo /usr/bin/rtorrent -n -o system.daemon.set=true ${opts} }
/usr/bin/rtorrent -n -o system.daemon.set=true ${opts}